# supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current/ext/CORS-Proxy
