# supervisord.conf
[program:node]
command = node ext/CORS-Proxy/server.js
directory = /home/dotcloud/current
