# supervisord.conf
[program:node]
command = node CORS-Proxy/server.js
directory = /home/dotcloud/current
